# systemd/personal-agent.service.template
[Unit]
Description=Personal Agent — FastAPI Backend (port {{API_PORT}})
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{ROOT_DIR}}/scripts
EnvironmentFile={{ROOT_DIR}}/scripts/backend/.env
ExecStart={{ROOT_DIR}}/scripts/backend/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port {{API_PORT}} --log-level info
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
